# amd.service
[Unit]
Description=AMD GPU CLOCK FIX

[Service]
ExecStart=/usr/bin/amdgpufixes

[Install]
WantedBy=multi-user.target 
